# srcs/Nginx/srcs/supervisord.conf
[supervisord]
user=root
nodaemon=true
[program:nginx]
command=nginx -g "daemon off;"